# utils/claude-keep-alive/com.terrancebrandon.claude-keep-alive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.terrancebrandon.claude-keep-alive</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/terrancebrandon/.local/bin/claude-keep-alive.sh</string>
    </array>

    <!-- Run every 30 minutes (1800 seconds) -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Run once immediately on load -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-keep-alive.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-keep-alive.err</string>
</dict>
</plist>
